# src/deb/control/postinst
#!/bin/sh

# update-rc.d sculptor-agent defaults
# update-rc.d sculptor-console defaults

# service sculptor-agent start
# service sculptor-console start