# src/deb/control/postrm
#!/bin/sh
update-rc.d sculptor-agent remove
rm -rf /var/galaxy