# src/deb/control/postinst
#!/bin/sh

update-rc.d sculptor-agent defaults
service sculptor-agent start