# src/deb/control/prerm
#!/bin/sh
# update-rc.d sculptor-agent remove
service sculptor-agent stop
