# src/deb/control/prerm
#!/bin/sh

update-rc.d -f sculptor-agent remove